# nfc_reader/supermarkt-nfc.service
[Unit]
Description=Kinder-Supermarkt NFC Reader Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/supermarket/nfc_reader
ExecStart=/home/pi/supermarket/nfc_reader/venv/bin/python reader.py --server http://supermarket-server.local:5000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
